# apple/OneLake/LaunchAgents/dev.debruyn.ofem.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Bundled LaunchAgent for the OFEM daemon.

  SMAppService (ServiceManagement) requires this plist to live at
  Contents/Library/LaunchAgents/<Label>.plist inside the app bundle.
  The Label key must match the plist filename (without .plist).

  The daemon binary lives at Contents/Helpers/ofem so it can be
  separately code-signed (daemon entitlements differ from the host app).
  We use the BundleProgram key (Apple's recommended form for SMAppService
  bundled agents) with a path relative to the bundle root. @executable_path
  is a dyld-runtime variable that launchd does NOT expand, despite what some
  online tutorials claim — using it causes launchd to look for a file
  literally named "@executable_path/../Helpers/ofem", which fails with
  EX_CONFIG (exit code 78). BundleProgram is resolved by launchd itself
  relative to the containing app bundle.

  BundleProgram and ProgramArguments are alternatives per launchd.plist(5)
  — combining them is non-canonical. BundleProgram alone is sufficient:
  launchd resolves the path relative to the bundle and does not need an
  explicit argv[0].

  Stdout/stderr are left unset intentionally: when launched by launchd
  via SMAppService the daemon writes to the unified log (os_log) and its
  log dir under the App Group container. If you need file-level logging
  during development add StandardOutPath/StandardErrorPath here.

  ProcessType Background tells launchd this is a long-lived helper, not
  an interactive process, which keeps it off the CPU-throttle front-end
  bucket. KeepAlive=true means launchd restarts the daemon if it exits
  unexpectedly (respawn). RunAtLoad=true means launchd starts it as soon
  as SMAppService.register() succeeds (i.e. at next login and immediately
  on the first registration call).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.debruyn.ofem.daemon</string>
    <key>BundleProgram</key>
    <string>Contents/Helpers/ofem</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>dev.debruyn.ofem</string>
    </array>
</dict>
</plist>
